# CloudStorage/scripts/systemd-timer-backup.service
[Unit]
Description="backup"

[Service]
ExecStart="gsutil cp -r /home/mpiase/data/ gs://testing-bucket-with-no-public-access/"